# lai-shi-lu.service
[Unit]
Description=交易来时路 Next.js App
After=network.target

[Service]
Type=simple
User=exedev
WorkingDirectory=/home/exedev/lai-shi-lu
EnvironmentFile=/home/exedev/lai-shi-lu/.env
ExecStart=/home/exedev/.local/bin/node /home/exedev/lai-shi-lu/node_modules/.bin/next start -p 3000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
